# systemd/system/smb.service.d/override.conf
[Unit]
Requires=mnt-storage3.mount
After=mnt-storage3.mount
