# systemd/system/iwd.service.d/override.conf
[Unit]
Requires=sys-subsystem-net-devices-wlp1s0.device
After=sys-subsystem-net-devices-wlp1s0.device

[Service]
ExecStart=
ExecStart=/usr/lib/iwd/iwd --nointerfaces 'wlan[0-9]'
