# systemd/system/custom-sleep-actions.service
[Unit]
Description=Customed sleep actions
Before=sleep.target
StopWhenUnneeded=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=-/usr/bin/systemctl stop yggdrasil.service
ExecStart=-/usr/bin/systemctl stop sing-box.service
ExecStop=-/usr/bin/iwctl station wlan0 scan
ExecStop=-/usr/bin/systemctl start yggdrasil.service
ExecStop=-/usr/bin/systemctl start sing-box.service

[Install]
WantedBy=sleep.target
