# systemd/system/lowbat.service
[Unit]
Description=Hibernate on low battery level

[Service]
Type=oneshot
ExecStart=/usr/bin/sh /etc/systemd/system/lowbat.sh
