# systemd/system/miner@.timer
[Unit]
Description=Run GPU miner for miner@%i.service

[Timer]
Unit=miner@%i.service
OnCalendar=*-*-* 21:00:00

[Install]
WantedBy=timers.target
