# systemd/system/miner@.timer
[Unit]
Description=Run GPU miner for miner@%i.service

[Timer]
OnCalendar=*-*-* 21:00:00
Unit=miner@%i.service

[Install]
WantedBy=timers.target
