# systemd/system/xmrig-stop.timer
[Unit]
Description=Stop Xmrig

[Timer]
OnCalendar=*-*-* 6:00:00

[Install]
WantedBy=timers.target
