# systemd/system/yggdrasil.service.d/override.conf
[Service]
StandardOutput=append:/tmp/yggdrasil.log
