# systemd/system/monerod.service.d/override.conf
[Unit]
Requires=mnt-storage1.mount
PartOf=mnt-storage1.mount
After=mnt-storage1.mount
